# capago-monitor.service
[Unit]
Description=Capago Alger visa appointment monitor (Long sejour -> Etudes)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/capago-monitor
ExecStart=/usr/bin/python3 %h/projects/capago-monitor/monitor.py
Restart=always
RestartSec=30
# Needed so notify-send / paplay / xdg-open reach your graphical session:
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=default.target
